# os/etc/cron.hourly/logrotate.sh___rm
#!/bin/bash

logrotate --skip-state-lock /etc/logrotate.d/app
